# supervisord.conf
; Supervisord config for Render: Streamlit + Flask + Caddy
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:streamlit]
command=%(ENV_ROOT)s/.venv-streamlit/bin/python -m streamlit run %(ENV_ROOT)s/app.py --server.port 8501 --server.address 0.0.0.0 --server.headless true --server.enableCORS false --server.enableXsrfProtection false --server.fileWatcherType none
directory=%(ENV_ROOT)s
autostart=true
autorestart=true
startsecs=10
startretries=100
stopwaitsecs=20
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
environment=STREAMLIT_BROWSER_GATHER_USAGE_STATS="false",STREAMLIT_SERVER_FILE_WATCHER_TYPE="none",PYTHONUNBUFFERED="1"

[program:flask]
command=%(ENV_ROOT)s/.venv-web/bin/gunicorn email_app:app --bind 127.0.0.1:5000 --timeout 180 --workers 1 --access-logfile - --error-logfile -
directory=%(ENV_ROOT)s
autostart=true
autorestart=true
startsecs=3
startretries=100
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1"

[program:caddy]
command=%(ENV_ROOT)s/bin/caddy run --config %(ENV_ROOT)s/Caddyfile --adapter caddyfile
directory=%(ENV_ROOT)s
autostart=true
autorestart=true
startsecs=2
startretries=100
redirect_stderr=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
priority=100
